# app/conf/supervisor_car.ini
[supervisord]
nodaemon=true

[program:car]
command=/app/car
user=root
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/app/logs/car_out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile=/app/logs/car_err.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB

[program:mongod]
command=/usr/bin/mongod --config /app/conf/mongod.conf
autostart=true
autorestart=true
user=mongodb
startsecs=10
startretries=3
stdout_logfile=/app/logs/mongo_out.log
stderr_logfile=/app/logs/mongo_err.log